# supportFiles/sysbackup.service

[Unit]
Description=Backing up drives

[Service]
Type=oneshot
ExecStart=echo "Test" > /tmp/test
